# templates/default/flink.service.erb
[Unit]
Description=Apache Flink
After=network.target

[Service]
Type=forking
LimitNOFILE=99999
ExecStart=<%= node['flink']['dir'] %>start-cluster.sh
ExecStop=<%= node['flink']['dir'] %>stop-cluster.sh
User=<%= node['flink']['user'] %>
Restart=on-failure
RestartSec=30s
[Install]
WantedBy=multi-user.target
