# templates/default/supervisor.conf.erb
[program:flink]
command=bash start-cluster.sh
directory=<%= node['flink']['dir'] %>
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/flink.log
